# last-shutdown-unexpected.target
[Unit]
Description=Last shutdown was unexpected
Requisite=last-shutdown-unexpected.service
After=last-shutdown-unexpected.service
